# RotateKeysDaemon/layout/Library/LaunchDaemons/de.tudarmstadt.seemoo.rotatekeys.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!-- Run time-based -->
<dict>
	<key>Disabled</key>
	<false/>
	<key>Label</key>
	<string>de.tudarmstadt.seemoo.rotatekeys</string>
	<key>OnDemand</key>
	<false/>
	<key>Program</key>
	<string>/usr/bin/rotatecgkeysd</string>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/dev/null</string>
</dict>
</plist>
